# deploy/firstboot/dune-firstboot.service
[Unit]
Description=dune first-boot per-instance credential generation
Documentation=https://github.com/leto217/DUNE
# Run after the network and cloud-init are up, but BEFORE the panel starts, so
# the panel never serves the default admin/admin account.
After=network-online.target cloud-init.service
Wants=network-online.target
Before=dune.service
# Skip entirely once the sentinel exists (cheap guard; the script re-checks too).
ConditionPathExists=!/etc/dune/.firstboot-done

[Service]
Type=oneshot
RemainAfterExit=yes
# Inherit the same DB configuration the panel uses (sqlite default / postgres).
EnvironmentFile=-/etc/default/dune
EnvironmentFile=-/etc/conf.d/dune
EnvironmentFile=-/etc/sysconfig/dune
ExecStart=/usr/local/dune/dune-firstboot.sh

[Install]
WantedBy=multi-user.target
